# roles/elasticsearch/files/elasticsearch.conf
[program:elasticsearch]
directory=/opt/elasticsearch
command=/opt/elasticsearch/bin/elasticsearch
autostart=true
autorestart=true
redirect_stderr=true
user=elasticsearch
